# services/transmission-daemon/run
#!/bin/sh
# Transmission daemon
exec /usr/bin/transmission-daemon -f --log-error
