# services/transmission-daemon/run
#!/bin/sh
exec chpst -u transmission:transmission /usr/bin/transmission-daemon -f --log-error
