# services/rsyslogd/run
#!/bin/sh
echo "Starting rsyslogd..."
exec rsyslogd -n
