# foo/etc/sv/agetty-generic/finish
#!/bin/sh
service=$(basename $(pwd))
tty=${service/*-}
exec utmpset -w $tty
