# services/dhcpcd/run
#!/bin/sh
exec /usr/sbin/dhcpcd -q -B 2>&1
