# services/dropbear/run
#!/bin/sh

echo "Starting dropbear..."
exec /usr/bin/dropbear -F -R
