# services/agetty-tty1/run
#!/bin/sh
tty=${PWD##*-}
exec setsid /sbin/agetty $([ $tty = tty1 ] && echo --noclear) $tty linux
