# services/agetty-tty1/run
#!/bin/sh
tty=${PWD##*-}
exec setsid /sbin/agetty -8 -s 38400 $([ $tty = tty1 ] && echo --noclear) $tty linux
