# services/agetty-ttyS0/run
#!/bin/sh
tty=${PWD##*-}
exec setsid /sbin/agetty -8 -s 9600 --noclear $tty linux
