# services/nfs-server/finish
#!/bin/sh

/usr/sbin/rpc.nfsd -- 0
/usr/sbin/exportfs -a -u
