# services/agetty-console/run
#!/bin/sh
if [ -x /sbin/getty ]; then
	GETTY=getty
elif [ -x /sbin/agetty ]; then
	GETTY=agetty
fi
exec $GETTY console 38400 linux
