# services/agetty-console/run
#!/bin/sh
exec /sbin/agetty console linux
