# services/xinetd/run
#!/bin/sh
exec /usr/sbin/xinetd -dontfork
