# services/upsdrvctl/run
#!/bin/sh
# Network UPS Tools - driver controller
exec /usr/bin/upsdrvctl -D start
