# services/alsa/run
#!/bin/sh
set -e
/usr/sbin/alsactl restore
exec chpst -b alsa pause
