# services/alsa/run
#!/bin/sh
set -e
/usr/bin/alsactl restore
exec chpst -b alsa pause
