# services/agetty-generic/run
#!/bin/sh
tty=${PWD##*-}
if [ -x /sbin/getty ]; then
	# busybox
	GETTY=getty
elif [ -x /sbin/agetty ]; then
	# util-linux
	GETTY=agetty
fi
exec setsid $GETTY $tty 38400 linux
