# services/agetty-generic/run
#!/bin/sh
tty=${PWD##*-}
exec setsid /sbin/agetty $tty linux
