# services/sulogin/run
#!/bin/sh
exec setsid /sbin/sulogin -p </dev/tty1 >/dev/tty1 2>&1
