# services/xbmc-standalone/run
#!/bin/sh
exec chpst -U xbmc:xbmc xbmc-standalone -l /run/lirc/lircd
