# services/apache/run
#!/bin/sh
exec /usr/sbin/httpd -DNO_DETACH
