# services/avahi-daemon/run
#!/bin/sh
sv check dbus >/dev/null || exit 1
echo "Starting avahi..."
exec avahi-daemon -s
