# services/ntpd/run
#!/bin/sh
echo "Starting ntpd..."
exec ntpd -g -u ntpd:ntpd -n
