# services/redis/run
#!/bin/sh
exec /usr/bin/redis-server /etc/redis/redis.conf
