# foo/etc/sv/dhcpcd/run
#!/bin/sh
echo "Starting dhcpcd..."
exec /usr/sbin/dhcpcd -q -B 2>&1
