# services/lighttpd/run
#!/bin/sh
exec /usr/sbin/lighttpd-angel -D -f /etc/lighttpd/lighttpd.conf
