# services/dcron/run
#!/bin/sh
echo "Starting dcron..."
exec crond -f 2>&1
